# docker-config/crontabs
#~ MAGENTO START 2d22aa451813cda4fbe2784e3e2f6f3a827532eea76ded46ef31a3e12155ea96
* * * * * /usr/bin/php8.1 /var/www/html/project-files/bin/magento cron:run 2>&1 | grep -v "Ran jobs by schedule" >> /var/www/html/project-files/var/log/magento.cron.log
#~ MAGENTO END 2d22aa451813cda4fbe2784e3e2f6f3a827532eea76ded46ef31a3e12155ea96